# bin/noip-updater
#!/bin/sh
### BEGIN INIT INFO
# Provides:          noip-update
# Required-Start:    $local_fs $remote_fs $network $syslog $named
# Required-Stop:     $local_fs $remote_fs $network $syslog $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     true
# Short-Description: NO-IP updater
# Description:       Start and stop the NO-IP updater
#  This script will start the service to update my NO-IP DNS entries with my public IP, if that public IP changes
### END INIT INFO

NODE_HOME=/home/noip/node-linux-armv7l

case "$1" in
  start)
     mv $LOG_DIR/log $LOG_DIR/log-$(date --iso-8601='seconds')
     su noip -c "cd /home/noip && $NODE_HOME/bin/node noip-updater/build/noip-updater.js 1>log 2>&1 &"
  ;;
  stop)
  ;;
  *)
    echo "Unknown argument: $1"
    exit 1;
  ;;
esac
